# Q2/2021213389苏鑫阳level4/sunshin.club(sh内所引用链接的file内容)/joshua_ddns.timer
[Unit]
Description=Joshua_DDNS

[Timer]
OnUnitActiveSec=1h 

[Install]
WantedBy=multi-user.target
