# Q2/2021214721张鑫辉/DDNS_Script/install_files/ddns.service
[Unit]
Description=DDNS Service
[Service]
Type=simple
ExecStart=/usr/local/ddns_script/ddns
