# Q2/2021213389苏鑫阳/sunshin.club(sh内所引用链接的file内容)/joshua_ddns.service
[Unit]
Description=Joshua_DDNS
Requires=network.target
[Service]
WorkingDirectory = /etc/joshua_ddns
ExecStart=
